# docker/puppetmaster/supervisord-puppet.conf
[program:puppet]
command=puppet master --no-daemonize
stdout_events_enabled = true
stderr_events_enabled = true

