# docker/puppetmaster/supervisord-foremanproxy.conf
[program:puppetdb]
user=foreman-proxy
command=/usr/share/foreman-proxy/bin/smart-proxy
environment=HOME=/usr/share/foreman-proxy
stdout_events_enabled = true
stderr_events_enabled = true

